# autossh/rootfs/etc/s6-overlay/s6-rc.d/autossh/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting autossh..."

exec /usr/bin/run.sh
